# systemd/nexus-exporter.service
[Unit]
Description=Nexus Exporter for Prometheus
Documentation=https://github.com/yimeng/nexus-exporter
After=network.target

[Service]
Type=simple
User=nexus-exporter
Group=nexus-exporter
WorkingDirectory=/var/lib/nexus-exporter

# Method 1: Using EnvironmentFile (variables loaded into environment)
EnvironmentFile=/etc/nexus-exporter/nexus-exporter.conf
ExecStart=/usr/local/bin/nexus-exporter

# Method 2: Alternative - use --config flag (uncomment if Method 1 doesn't work)
# ExecStart=/usr/local/bin/nexus-exporter --config=/etc/nexus-exporter/nexus-exporter.env

Restart=always
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/nexus-exporter

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
